# service/cache/run
#!/usr/bin/env bash
set -e
exec 2>&1

export API_BASE_URL="${API_BASE_URL:-http://arid:3000}"

cd "/app" || exit 1
exec node "/app/server/index.mjs"